# FSW/crontab
@reboot sleep 10 && /home/pi/scripts/RunMain.sh
